# rag_api/rootfs/etc/services.d/postgres/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

DATA_DIR="/data/postgres"

# Initialize PostgreSQL if it hasn't been initialized
if [ ! -d "$DATA_DIR" ]; then
    bashio::log.info "Initializing PostgreSQL database..."
    mkdir -p "$DATA_DIR"
    chown postgres:postgres "$DATA_DIR"
    su postgres -c "initdb -D $DATA_DIR"

    # Create database and enable pgvector extension
    su postgres -c "pg_ctl -D $DATA_DIR -l /var/log/postgresql/postgresql.log start"
    su postgres -c "createdb ragapi"
    su postgres -c "psql -d ragapi -c \"CREATE EXTENSION IF NOT EXISTS vector;\""
    su postgres -c "pg_ctl -D $DATA_DIR stop"
fi

bashio::log.info "Starting PostgreSQL..."
exec s6-setuidgid postgres postgres -D "$DATA_DIR"
