# rag_api/rootfs/etc/services.d/postgres/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

# Check if PostgreSQL data directory is empty
if [ -z "$(ls -A "$PGDATA")" ]; then
    bashio::log.info "Initializing PostgreSQL database..."
    su postgres -c "initdb -D $PGDATA"

    # Modify postgresql.conf to enable logging
    echo "log_statement = 'all'" >> "$PGDATA/postgresql.conf"
    echo "log_destination = 'stderr'" >> "$PGDATA/postgresql.conf"

    # Start PostgreSQL temporarily to create database and extension
    su postgres -c "pg_ctl -D $PGDATA -l /var/log/postgresql/postgresql.log start"

    # Wait for PostgreSQL to start
    until su postgres -c "pg_isready" 2>/dev/null; do
        bashio::log.info "Waiting for PostgreSQL to start..."
        sleep 1
    done

    # Set a password for postgres user
    su postgres -c "psql -c \"ALTER USER postgres WITH PASSWORD 'ragapipassword';\""

    # Create database and enable pgvector extension
    su postgres -c "createdb ragapi"
    su postgres -c "psql -d ragapi -c 'CREATE EXTENSION IF NOT EXISTS vector;'"

    # Stop PostgreSQL
    su postgres -c "pg_ctl -D $PGDATA stop"
fi

bashio::log.info "Starting PostgreSQL..."
exec su postgres -c "postgres -D $PGDATA"
