# librechat/rootfs/etc/services.d/meilisearch/finish
#!/usr/bin/env bashio
# ==============================================================================
# Handle Meilisearch service termination
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

declare APP_EXIT_CODE=${1}

if [[ "${APP_EXIT_CODE}" -ne 0 ]] && [[ "${APP_EXIT_CODE}" -ne 256 ]]; then
  # Check if this is an address-in-use error
  if [[ -f /tmp/meilisearch_output ]] && grep -q "Address already in use" /tmp/meilisearch_output; then
    bashio::log.warning "Meilisearch failed to start: Address already in use. Search functionality will be disabled."
  else
    bashio::log.warning "Meilisearch terminated with exit code ${APP_EXIT_CODE}. Search functionality will be disabled."
  fi
  # Don't halt the entire add-on, LibreChat can function without Meilisearch
fi

bashio::log.info "Service restart after closing"
