# librechat/rootfs/etc/services.d/meilisearch/run
#!/usr/bin/with-contenv bashio
## Run meilisearch
bashio::log.info "Starting Meilisearch"
exec /opt/meilisearch
