# librechat/rootfs/etc/services.d/meilisearch/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
## Run meilisearch
bashio::log.info "Starting Meilisearch"
mkdir -p /data/meilisearch
exec /opt/meilisearch --db-path /data/meilisearch
