# librechat/rootfs/etc/services.d/librechat/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
cd /opt/librechat || exit
bashio::log.info "Starting LibreChat Backend"

#if there isn't a default env file, then copy one.
if [ ! -f /config/.env ]; then
    cp .env.example /config/.env
fi

#copy whatever is in env before running.
cp /config/.env .env

exec npm run backend
