# librechat/rootfs/etc/services.d/mongodb/finish
#!/usr/bin/env bashio
if [[ "${1}" -ne 0 ]]; then
  bashio::log.warning "MongoDB crashed, halting add-on"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "MongoDB stopped, restarting..."