# librechat/rootfs/etc/services.d/mongodb/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
## Run your program
bashio::log.info "Starting MongoDB"
mkdir -p /data/db
exec mongod --bind_ip 0.0.0.0 